# bin/gitsync-dirs.target
[Unit]
Description=Automatically will do git-sync on repos

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
